# backend/flask_supervisor.conf
[program:flask_nginx]
command = gunicorn wsgi -b localhost:5000
directory = /app
autostart = true
autorestart = true
